# templates/systemd/certbot-renew.service.j2
# {{ ansible_managed }}

[Unit]
Description=Renew Certbot certificates

[Service]
Type=oneshot
ExecStart={{ funkwhale_letsencrypt_venv_dir }}/bin/certbot renew --quiet
