# templates/systemd/certbot-renew.timer.j2
# {{ ansible_managed }}

[Unit]
Description=Weekly Certbot certificates renewal

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
